# scripts/proxmoxvue-agent.service
[Unit]
Description=ProxmoxVue Agent
Documentation=https://github.com/TheLion/proxmoxvue-agent
After=network-online.target pve-cluster.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/proxmoxvue-agent run
Restart=on-failure
RestartSec=5
User=root
Group=root
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# /var/log staat expliciet ReadWrite zodat lumberjack de log-file en zijn
# rotated backups (.1, .2, .gz) in /var/log/ kan schrijven. Nodig omdat
# ProtectSystem=strict heel /var/log read-only zou maken zonder deze regel.
ReadWritePaths=/etc/proxmoxvue-agent /var/log
# Cache voor de server-driven config (proxmoxvue.app/config/v1.json).
# StateDirectory creëert /var/lib/proxmoxvue-agent met owner=User automatisch
# en voegt het impliciet aan ReadWritePaths toe. Persistent over restarts.
StateDirectory=proxmoxvue-agent
PrivateTmp=true
# Single source of truth: agent schrijft direct naar /var/log/proxmoxvue-agent.log
# via lumberjack. Stdout heeft niets meer; null voorkomt dat journald lege
# regels accumuleert. Stderr blijft naar journal — daar landen panics en
# pre-slog fatals.
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
